# cookbooks/superdesk/templates/default/superdesk.conf.erb

description "Superdesk"

start on startup
stop on shutdown

respawn

exec python3 <%= @app %> 2>&1 >> /tmp/superdesk.log
