# cookbooks/superdesk/templates/default/superdesk.conf.erb

description "Superdesk"

start on startup
stop on shutdown

respawn

env PYTHONPATH="<%= @pythonpath %>"

exec python3 <%= @app %> 2>&1 >> /tmp/superdesk.log
