# cookbooks/superdesk-web/templates/default/superdesk-web.conf.erb

description "Superdesk Web"

start on startup
stop on shutdown

respawn

env MONGOHQ_URL="<%= @db %>"
exec python <%= @app %> runserver 0.0.0.0:8000 2>&1 >> /tmp/superdesk-web.log
